# examples/systemd/net-watchdog.timer
[Unit]
Description=Run Net Watchdog periodically

[Timer]
OnBootSec=30s
OnUnitActiveSec=20s
AccuracySec=1s

[Install]
WantedBy=timers.target
